# w_level_5.conf
[program:w_level_5_back]
command=sudo python /home/pi/w_level_5.py
numprocs=1
redirect_stderr=true
stdout_logfile=/var/log/w_level_5.log
user=root

